# deploy/cron/laravel-scheduler.cron
* * * * * cd /var/www/html && php artisan schedule:run >> storage/logs/scheduler.log 2>&1
